# examples/systemd/portal-hub-web.service
[Unit]
Description=Portal Hub web API and OAuth service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=portal-hub
Group=portal-hub
Environment=PORTAL_HUB_STATE_DIR=/var/lib/portal-hub
Environment=PORTAL_HUB_SSH_PORT=2222
# Set this when using Tailscale Serve or another HTTPS reverse proxy.
# Environment=PORTAL_HUB_PUBLIC_URL=https://portal-hub.example.ts.net
ExecStart=/usr/local/bin/portal-hub web --bind 127.0.0.1:8080
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
